# roles/unicorn/templates/unicorn.service.j2
[Unit]
SourcePath=/etc/init.d/unicorn_{{ app }}
[Service]
ExecStart=/etc/init.d/unicorn_{{ app }} start
ExecStop=/etc/init.d/unicorn_{{ app }} stop
ExecReload=/etc/init.d/unicorn_{{ app }} reload
